# scripts/github-activity.timer
[Unit]
Description=Daily GitHub activity fetch for manueloelmaier.de

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
