# scripts/github-activity.service
[Unit]
Description=Fetch GitHub activity for manueloelmaier.de widget
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Adjust path to match your deployment directory on the server
ExecStart=/var/www/manueloelmaier.de/scripts/fetch-github.sh
StandardOutput=journal
StandardError=journal
